# services.d/template/run
#!/usr/bin/execlineb -P
sh -c /etc/consul.startup